# src/bonjour-daemon/linux-bonjour.service
[Unit]
Description=Linux Bonjour Biometric Authentication Daemon
After=network.target network-online.target camera.target
Wants=camera.target

[Service]
Type=simple
ExecStart=/usr/bin/linux-bonjour-daemon
Restart=always
RestartSec=3
# Ensure the runtime directory exists for the UDS socket
RuntimeDirectory=linux-bonjour
# Ensure the state directory exists for biometric signatures
StateDirectory=linux-bonjour
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
