# config/linux-bonjour.service
[Unit]
Description=Linux Bonjour Biometric Authentication Daemon
After=network.target network-online.target camera.target
Wants=camera.target
Before=display-manager.service

[Service]
Type=notify
ExecStart=/usr/bin/linux-bonjour-daemon
Restart=always
RestartSec=3
Environment=RUST_LOG=info

# Standard Runtime Directories
RuntimeDirectory=linux-bonjour
StateDirectory=linux-bonjour

# Hardening and Sandboxing
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE
SupplementaryGroups=tss
DeviceAllow=/dev/tpm0 rw
DeviceAllow=/dev/tpmrm0 rw
DeviceAllow=/dev/video0 rw
DeviceAllow=/dev/video1 rw
DeviceAllow=/dev/video2 rw
DeviceAllow=/dev/video3 rw
DeviceAllow=/dev/video4 rw
DeviceAllow=/dev/video5 rw
DeviceAllow=/dev/video6 rw

[Install]
WantedBy=multi-user.target
