# config/linux-bonjour.service
[Unit]
Description=Linux Bonjour Face Recognition Daemon
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/share/linux-bonjour
ExecStart=/usr/share/linux-bonjour/venv/bin/python src/daemon/main.py
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
